# provision/roles/radius/files/docker/files/etc/services.d/radiusd/run
#!/usr/bin/execlineb -P

with-contenv
/usr/sbin/radiusd -f -l stdout
