# provision/roles/radius/files/docker/files/etc/services.d/winbindd/run
#!/usr/bin/execlineb -P

/usr/sbin/winbindd --foreground --stdout
